# target/MyRetail/build/resources/main/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=mongod

[program:apache2]
command=/usr/bin/java -cp /tmp/app org.springframework.boot.loader.JarLauncher